# com.local.xplock-reinject.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.xplock-reinject</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/draistrick/.local/bin/xplock-reinject</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/draistrick/Library/Logs/xplock-reinject.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/draistrick/Library/Logs/xplock-reinject.err</string>
</dict>
</plist>
