# scripts/pkg-scripts/postinstall
#!/bin/bash
# postinstall: bootstrap LaunchAgent + restart Finder if user agreed preinstall.
set -e

BIN="/Library/Application Support/ColumnTamer/ColumnTamerHelper"
HELPER_PLIST="/Library/LaunchAgents/columntamer.helper.plist"
MENU_PLIST="/Library/LaunchAgents/columntamer.menu.plist"
FLAG="/var/run/.columntamer.restart"

chmod 755 "$BIN"

# M7: remove stale shared logdir from older installs (was 1777 symlink trap).
# Logs now per-user at ~/Library/Logs/ColumnTamer/.
rm -rf "/Library/Application Support/ColumnTamer/logs"

# M6: pin ownership to root:wheel (installer runs as root; staged files
# otherwise may carry dev-user ownership = persistence vector).
chown -R root:wheel "/Library/ScriptingAdditions/ColumnTamer.osax"
chown -R root:wheel "/Library/Application Support/ColumnTamer"
chown root:wheel "$HELPER_PLIST" "$MENU_PLIST"

CONSOLE_USER="$(/usr/bin/stat -f%Su /dev/console 2>/dev/null || true)"
CONSOLE_UID="$(/usr/bin/id -u "$CONSOLE_USER" 2>/dev/null || true)"

# bootstrap agents AS USER. RunAtLoad=true auto-launches; no kickstart (blocks).
# bootstrap runs backgrounded so postinstall returns fast.
# NOTE: postinstall run as root by package_script_service — no tty.
# sudo -u fail here. root can bootstrap user gui domain directly.
if [[ -n "$CONSOLE_UID" ]]; then
  # NOTE: launchctl bootstrap gui/$UID fail from package_script_service
  # (EIO sandbox restriction). Plists install correctly — launchd picks
  # them up on next login via RunAtLoad. For immediate effect, launch
  # binaries directly via launchctl asuser (spawn in user context).
  /usr/bin/killall ColumnTamerMenu 2>/dev/null || true
  launchctl asuser "$CONSOLE_UID" /usr/bin/open "/Library/Application Support/ColumnTamer/ColumnTamerMenu.app" 2>/dev/null &
  if [[ -x "$BIN" ]]; then
    launchctl asuser "$CONSOLE_UID" "$BIN" 2>/dev/null &
  fi
fi

# restart only if user agreed at preinstall prompt
if [[ -f "$FLAG" ]]; then
  /usr/bin/killall Finder 2>/dev/null || true
  rm -f "$FLAG"
else
  if [[ -n "$CONSOLE_UID" ]]; then
    launchctl asuser "$CONSOLE_UID" /usr/bin/osascript -e \
      'display notification "ColumnTamer installed. Restart Finder to activate." with title "ColumnTamer"' \
      2>/dev/null || true
  fi
fi

exit 0
